# .github/actions/postgresql/postgresql.conf
[program:postgresql]
command=/etc/init.d/postgresql start
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/postgresql.err.log
stderr_logfile=/var/log/supervisor/postgresql.out.log
